# visual/themes/themes/GTK/external/aerothemeplasma/misc/libplasma/apply-libplasma-patches.service
[Unit]
Description=Apply AeroThemePlasma libplasma patches
Before=
[Service]
Type=simple
ExecStart=/opt/aerothemeplasma/tmp/apply
[Install]
WantedBy=multi-user.target
